# etc/rc.local
#!/bin/sh

/opt/emailnotify.sh &

exit 0
